# Infrastructure/node-red.service
[Unit]
Description=Node-RED - Flow-based programming for the Internet of Things
Documentation=http://nodered.org/docs/
After=syslog.target network.target redis.service

[Service]
Type=simple
User=antoine
Group=antoine
WorkingDirectory=/home/antoine/.node-red
ExecStart=/usr/local/bin/node-red --settings /home/antoine/.node-red/settings.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

